# ansible/roles/quickfeed/templates/systemd/quickfeed.service.d/agenv.conf.j2
[Service]
Environment="GITHUB_KEY={{ quickfeed_github_key }}"
Environment="GITHUB_SECRET={{ quickfeed_github_secret }}"
